# etc/supervisor/conf.d/mastodon.conf
[program:mastodon]
command=/usr/local/sbin/start-mastodon

user=mastodon

priority=20

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
